# .docker/build/php/nginx-supervisor.ini
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true

[program:php-fpm]
command = /usr/sbin/php5-fpm -F
user = root
autostart = true

[program:mysql]
command = /run.sh
user = root
autostart = true

[program:permissions]
command = /opt/scripts/permissions.sh
user = root
autostart = true
